# scripts/systemd/uqm-watchdog.timer
[Unit]
Description=UQM Health Watchdog Timer (every 2 min)

[Timer]
OnBootSec=90
OnUnitActiveSec=2min
AccuracySec=30
Persistent=true

[Install]
WantedBy=timers.target default.target
